# docker/supervisor/conf/supervisord-worker.conf
[program:app_worker]
command=/usr/bin/php /app/worker.php
autorestart=true
stdout_logfile=/proc/1/fd/1
redirect_stderr=true

